# trunk/target/raspberrypi.eglibc/rootfs/etc/inittab
::sysinit:/sbin/syslogd -S -C512
::sysinit:/etc/rc.bootup
ttyAMA0::respawn:/etc/rc.initial

