# target/beaglebone-black.eglibc/rootfs/etc/inittab
::sysinit:/sbin/syslogd -S -C512
::sysinit:/etc/rc.bootup.php
ttyO0::respawn:/etc/rc.initial

